# deployment/runner_supervisor.conf
[program:sharpie-runner]
directory=/var/www/sharpie/runner/

# Activate your virtual environment and start the runner
# IMPORTANT: Replace YOUR_CONNECTION_KEY with a key created in Django Admin
# See SERVER_SETUP.md Step 11.5 for instructions
environment=PATH=/var/www/sharpie/venv/bin/
command=sharpie-runner runserver --connection-key=YOUR_CONNECTION_KEY

# Log file location
stdout_logfile=/var/www/sharpie/logs/runner.log
redirect_stderr=true

# Automatically start and recover processes
autostart=true
autorestart=true